# scripts/rpi/tecpan-mobile.service
[Unit]
Description=Tecpancito POS Mobile (puerto 5174)
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__/apps/pos-mobile
ExecStart=__REPO_DIR__/node_modules/.bin/vite preview --port 5174 --host 0.0.0.0
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tecpan-mobile

LimitNOFILE=65536
MemoryMax=128M

[Install]
WantedBy=multi-user.target
